# launchd/com.nikhilesh.leadtracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nikhilesh.leadtracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd /Users/nikhileshgupta/shree-madhav-lead-tracker &amp;&amp; /Users/nikhileshgupta/shree-madhav-lead-tracker/.venv/bin/python -m lead_tracker daily --input-dir /Users/nikhileshgupta/whatsapp_exports</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/nikhileshgupta/Library/Logs/leadtracker.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nikhileshgupta/Library/Logs/leadtracker.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>TZ</key>
        <string>Asia/Kolkata</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
